# pi/systemd/rfc2217-hotplug@.service
[Unit]
Description=RFC2217 Hotplug Handler for %i
After=network.target rfc2217-portal.service
Wants=rfc2217-portal.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rfc2217-hotplug add /dev/%i
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal
