# pi/udev/99-rfc2217.rules
# RFC2217 hotplug rules for USB serial devices
# Uses systemd-run to properly escape udev environment

# ttyACM devices (CDC-ACM: ESP32, Arduino, etc.)
ACTION=="add", SUBSYSTEM=="tty", KERNEL=="ttyACM*", RUN+="/usr/bin/systemd-run --no-block /usr/local/bin/rfc2217-hotplug add /dev/%k %E{ID_PATH}"
ACTION=="remove", SUBSYSTEM=="tty", KERNEL=="ttyACM*", RUN+="/usr/local/bin/rfc2217-hotplug remove /dev/%k %E{ID_PATH}"

# ttyUSB devices (FTDI, CP210x, CH340, etc.)
ACTION=="add", SUBSYSTEM=="tty", KERNEL=="ttyUSB*", RUN+="/usr/bin/systemd-run --no-block /usr/local/bin/rfc2217-hotplug add /dev/%k %E{ID_PATH}"
ACTION=="remove", SUBSYSTEM=="tty", KERNEL=="ttyUSB*", RUN+="/usr/local/bin/rfc2217-hotplug remove /dev/%k %E{ID_PATH}"
